# ops/backup/dataspace-postgres-backup.timer
[Unit]
Description=Run Dataspace PostgreSQL backup evidence job daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=dataspace-postgres-backup.service

[Install]
WantedBy=timers.target
